# cron.d/sidekiq-cron
*/1 * * * * mastodon cd /home/mastodon/mastodon; RAILS_ENV=production /home/mastodon/.rbenv/shims/bundle exec rails runner 's = Sidekiq::Stats.new; print "processed:", s.processed, "\n"; print "failed:", s.failed, "\n"; print "default:", s.queues["default"], "\n"; print "mailers:", s.queues["mailers"], "\n"; print "push:", s.queues["push"], "\n"; print "pull:", s.queues["pull"], "\n"; print "latency:", Sidekiq::Queue.new.latency, "\n"' > /var/tmp/sidekiq-info
